# .config/systemd/user/aw-afk.service
[Unit]
Description=ActivityWatch Watcher: AFK

[Service]
ExecStart=/sbin/aw-watcher-afk
RestartSec=30
Restart=on-failure


[Install]
WantedBy=default.target
