# .config/systemd/user/aw-afk.service
[Unit]
Description=ActivityWatch Watcher: AFK

[Service]
ExecStart=/sbin/aw-watcher-afk


[Install]
WantedBy=default.target
