# .config/systemd/user/aw-server.service
[Unit]
Description=ActivityWatch server

[Service]
ExecStart=/sbin/aw-server


[Install]
WantedBy=default.target
