# .config/systemd/user/aw-window.service
[Unit]
Description=ActivityWatch Watcher: Window

[Service]
ExecStart=/sbin/aw-watcher-window


[Install]
WantedBy=default.target
